# cron/crontab
# This is an example cron job that runs every minute.
# * * * * * /path/to/script.sh
* * * * * /root/script.sh
