# cron/crontab
# This is an example cron job that runs every minute.
# * * * * * /path/to/script.sh
00 01,13 * * * /root/script.sh # it will run at 1 A.M and 1 P.M
